# ansible/roles/mikrotik_integration/templates/mikrotik-healthcheck.service.j2
# -----------------------------------------------------------------------------
# Filename: ansible/roles/mikrotik_integration/templates/mikrotik-healthcheck.service.j2
# Created: 2026-01-19
# Description: Systemd service to run MikroTik health checks.
# Developer notes:
#   - Health checks are non-invasive and safe to run frequently.
# ----------------------------------------------------------------------------

[Unit]
Description=MikroTik healthcheck (homelab_2026.2)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
EnvironmentFile=-/root/.config/homelab_2026_2/state.env
ExecStart=/opt/homelab/mikrotik/healthcheck.sh
